# script/supervisor/inowas.rabbit.flopy.runner.conf
[program:inowas.flopy.rabbit.runner]
command=php /var/www/modflowMS/bin/console inowas:flopy:rabbit:runner
